# scripts/vex-api.service
[Unit]
Description=Vex API Server
After=network.target ollama.service

[Service]
Type=simple
User=vex
Group=vex
WorkingDirectory=/opt/vex-api
ExecStartPre=-+/opt/vex-api/scripts/vex-update.sh
ExecStart=/opt/vex-api/vex-api
Restart=always
RestartSec=5

# Environment (ignore if missing)
EnvironmentFile=-/opt/vex-api/.env

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
