# scripts/vex-update.service
[Unit]
Description=Update Vex compiler to latest release

[Service]
Type=oneshot
ExecStart=/opt/vex-api/scripts/vex-update.sh
Environment=VEX_DIR=/opt/vex
Environment=VEX_REPO=meftunca/vex
TimeoutStartSec=180
